# scripts/etc/rc.d/dnsmasq
#!/bin/sh
#

daemon="/usr/local/sbin/dnsmasq"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
